# etc/rc.local
#!/bin/sh

# Create a directory for Tor control socket
mkdir -m 700 -p /run/tor
chown tor:tor /run/tor

# ...a separate mountpoint for ovpn config uploads
mkdir -p /mnt/phpd_upload
mount -o noexec,nosuid,nodev -t tmpfs tmpfs /mnt/phpd_upload

# copy firmware release information to phpd-readable location
cp nak-timestamp nak-release /tmp
